# gunicorn_app/templates/upstart.conf.erb
description "upstart configuration for gunicorn app <%= name %>"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec sudo mkdir -p /var/run/gunicorn
exec sudo chown www-data:www-data /var/run/gunicorn
exec sudo -u www-data <%= gunicorn_bin %> -c <%= config_path %> <%= settings_path %>